# raspberrypi-ua-netinst/config/files/root/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin screen --noclear %I 38400 linux
